# systemd/gdrive-autosync.service
[Unit]
Description=Google Drive auto-sync (inotify + rclone, one-way local to Drive)
After=network-online.target

[Service]
ExecStart=%h/.local/bin/gdrive-autosync daemon
Restart=on-failure
RestartSec=30
Nice=10

[Install]
WantedBy=default.target
