# systemd/gsync-mount.service
[Unit]
Description=gsync — browse Google Drive at ~/GoogleDrive (rclone mount)
Documentation=https://github.com/oguzkaganozt/gsync
After=network-online.target

[Service]
Type=notify
ExecStartPre=/bin/mkdir -p %h/GoogleDrive
ExecStart=/usr/bin/rclone mount gdrive: %h/GoogleDrive \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 2G \
    --dir-cache-time 1m \
    --poll-interval 30s
ExecStop=/bin/sh -c 'fusermount -uz %h/GoogleDrive 2>/dev/null || fusermount3 -uz %h/GoogleDrive'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
