# systemd/gsync-tray.service
[Unit]
Description=gsync system tray indicator
Documentation=https://github.com/oguzkaganozt/gsync
PartOf=graphical-session.target
After=graphical-session.target gsync.service

[Service]
Type=simple
ExecStart=%h/.local/bin/gsync-tray
Restart=on-failure
RestartSec=3
Environment=XDG_RUNTIME_DIR=%t

[Install]
WantedBy=graphical-session.target
